# autoauth
#!/bin/sh /etc/rc.common

START=99

start() {
    PID=$(ps | grep auto-auth:S | grep -v grep | awk '{print $1}')
    if [ ! -n "$PID" ]; then
        /etc/autoauth.sh 190000000000 1234 auto-auth:S &
    else
        echo The program is running, please do not run it repeatedly!
    fi
}
stop() {
    PID=$(ps | grep auto-auth:S | grep -v grep | awk '{print $1}')
    if [ -n "$PID" ]; then
        kill -9 $PID
    fi
}
restart() {
    stop
    start
}
